# backend/s6/s6-rc.d/redis/run
#!/command/with-contenv sh
exec redis-server --bind 127.0.0.1 --port 6379 --dir /data/redis --save 60 1000 --appendonly no --loglevel notice
